# contrib/checkdiff.service
[Unit]
Description=checkdiff daemon
# Network is needed for fetching source URLs and publishing to
# ntfy. After=network-online.target waits for DHCP / link-local
# addressing to complete so the first check doesn't race the
# network coming up. These are user-instance targets; for a
# system-wide unit, change to network-online.target and
# WantedBy=multi-user.target below.
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# %h expands to the user's home directory in the user systemd
# instance. -config/-state default to the user's
# ~/.config/checkdiff/config.toml and ~/.local/share/checkdiff/
# state.json, so they could be omitted; spelling them out
# matches what the README documents.
ExecStart=%h/bin/checkdiff -config %h/.config/checkdiff/config.toml -state %h/.local/share/checkdiff/state.json -v
WorkingDirectory=%h
Restart=always
RestartSec=5

# checkdiff writes only to the state file (atomically, via
# tmp-file + rename). The user systemd instance already gives
# the service access to the user's home; ProtectSystem=strict
# makes everything outside the user's home read-only as a
# defence-in-depth measure. We don't set ProtectHome=read-only
# because the daemon needs to write to
# ~/.local/share/checkdiff/state.json.
ProtectSystem=strict
PrivateTmp=true

[Install]
WantedBy=default.target
