# contrib/checkdiff.service
[Unit]
Description=checkdiff daemon
# Network is needed for fetching source URLs and publishing to
# ntfy. After=network-online.target waits for DHCP / link-local
# addressing to complete so the first check doesn't race the
# network coming up.
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/checkdiff -config /etc/checkdiff/config.toml -state /var/lib/checkdiff/state.json
Restart=always
RestartSec=5

# The daemon writes only to the state file. StateDirectory=
# makes /var/lib/checkdiff writable and owned by the service
# user; the binary is told to put the state file there via
# the -state flag. ProtectSystem=strict makes everything else
# read-only.
StateDirectory=checkdiff
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true

[Install]
WantedBy=multi-user.target
